# control-plane/systemd/ema-detector.service
[Unit]
Description=EMA/Z-Score Anomaly Detector — time-series baseline comparison
After=network-online.target postgresql.service ingestion.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/control-plane/detector
EnvironmentFile=/opt/control-plane/detector/.env
ExecStart=/opt/control-plane/venv/bin/python /opt/control-plane/detector/ema_detector.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ema-detector

[Install]
WantedBy=multi-user.target
