# control-plane/systemd/health-failover.service
[Unit]
Description=Health Check Endpoint - Failover Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/health
ExecStart=/opt/health/venv/bin/gunicorn \
    --bind 0.0.0.0:8080 \
    --workers 2 \
    --timeout 30 \
    --access-logfile - \
    app:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
